# Oracle/luigi/supervisor_luigi.ini
[supervisord]
nodaemon=true

[program:luigid]
autostart = true
autorestart = true
command=luigid
stderr_logfile=/var/log/luigid.err.log
stdout_logfile=/var/log/luigid.out.log
