# Oracle/airflow/supervisor_airflow.ini
[supervisord]
nodaemon=true

[program:airflow_web]
autostart = true
autorestart = true
command=airflow webserver -p 8080
stderr_logfile=/var/log/airflow_web.err.log
stdout_logfile=/var/log/airflow_web.out.log


[program:airflow_scheduler]
autostart = true
autorestart = true
command=airflow scheduler
stderr_logfile=/var/log/airflow_scheduler.err.log
stdout_logfile=/var/log/airflow_scheduler.out.log
